# templates/ansible/roles/hermes/templates/hermes-vps-healthcheck.timer.j2
[Unit]
Description=Run Hermes VPS healthcheck timer

[Timer]
OnCalendar={{ hermes_healthcheck_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
