# templates/ansible/roles/vps_ops/templates/hermes-vps-healthcheck.service.j2
[Unit]
Description=Run Hermes VPS healthcheck
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hermes-vps-healthcheck
