# templates/ansible/roles/vps_ops/templates/hermes-docker-cleanup.service.j2
[Unit]
Description=Prune unused Docker artifacts for Hermes VPS

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hermes-docker-cleanup
