# templates/ansible/roles/vps_ops/templates/hermes-vps-backup.timer.j2
[Unit]
Description=Run Hermes VPS backup

[Timer]
OnCalendar={{ hermes_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
