# templates/ansible/roles/vps_ops/templates/hermes-vps-release-check.timer.j2
[Unit]
Description=Run Hermes release check

[Timer]
OnCalendar={{ hermes_release_check_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
