# templates/ansible/roles/vps_ops/templates/hermes-vps-backup.service.j2
[Unit]
Description=Create Hermes VPS backup

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hermes-vps-backup
