# templates/ansible/roles/n8n/templates/hermes-n8n-sqlite-backup.service.j2
[Unit]
Description=Hermes n8n SQLite online backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=root
ExecStart={{ n8n_sqlite_backup_script_path }}
