# templates/ansible/roles/n8n/templates/hermes-n8n-sqlite-backup.timer.j2
[Unit]
Description=Run Hermes n8n SQLite online backup

[Timer]
OnCalendar={{ n8n_sqlite_backup_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
